# templates/default/janusgraph.service.erb
[Unit]
Description = janusgraph
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node['janusgraph']['user'] %>

PIDFile = <%= node['janusgraph']['pid_file'] %>
ExecStart = "<%= node['janusgraph']['base_dir'] %>/bin/janusgraph.sh start"
ExecStop = "<%= node['janusgraph']['base_dir'] %>/bin/janusgraph.sh stop"
LimitNOFILE = 20000

Type = forking

[Install]
WantedBy = multi-user.target
